# scripts/com.amass.leader-wiki-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amass.leader-wiki-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO_ROOT__/scripts/sync-wiki.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>WIKI_LOCAL_PATH_HOST</key>
        <string>__WIKI_LOCAL_PATH_HOST__</string>
        <key>SYNC_LOG_FILE</key>
        <string>__HOME__/Library/Logs/leader-wiki-sync.log</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/leader-wiki-sync.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/leader-wiki-sync.log</string>
</dict>
</plist>
